# mongodb3/templates/default/mongod.upstart.erb
# Generated by Chef - manual changes made will be lost!
description     "MongoDB server instance"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# If the process quits unexpectedly trigger a respawn
respawn

# Increase ulimit to recommended levels
limit nofile <%= node['mongodb3']['mongod']['open_file_limit'] %> <%= node['mongodb3']['mongod']['open_file_limit'] %>
limit nproc  <%= node['mongodb3']['mongod']['max_proc_number'] %> <%= node['mongodb3']['mongod']['max_proc_number'] %>

# Start the process
exec start-stop-daemon --start --chuid <%= node['mongodb3']['user'] %> --make-pidfile --pidfile <%= @pid_file %> --exec /usr/bin/mongod -- -f <%= @config_file %>

# Clean up after stopping the process
post-stop script
        rm <%= @pid_file %>
end script
